# puppypi-device-connect.service.example
[Unit]
Description=Device Connect read-only PuppyPi driver
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/device_connect
Environment=DEVICE_PROFILE=puppypi
ExecStart=/home/pi/device_connect/puppy-device-connect/start_fabric.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
