# puppypi-device-connect.service.example
[Unit]
Description=Device Connect read-only PuppyPi driver
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/odin/device-connect
Environment=PROJECT_ROOT=/home/odin/device-connect
Environment=MESSAGING_CREDENTIALS_FILE=/home/odin/device-connect/beta-odin-pi5.creds.json
ExecStart=/home/odin/device-connect/puppy-device-connect/start_fabric.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
